# deploy/usparts-healthcheck.timer
[Unit]
Description=Run USParts health check every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
